# file_compression.service
[Unit]

[Service]
ExecStart=/usr/bin/python3 /root/files_compression/app.py

[Install]
