# file_compression.service
[Unit]

[Service]
ExecStart=/root/files_compression/app.py

[Install]
